# etc/cron.daily/empty-trash
#!/usr/bin/bash

rm -rf /home/james/.local/share/Trash/files/* /home/james/.local/share/Trash/info/*
